# deploy/opsfleet-backend.service.example
[Unit]
Description=OpsFleetPilot API (ft-backend)
After=network.target postgresql.service redis-server.service

[Service]
Type=simple
User=root
WorkingDirectory=/OPSFLEET_ROOT/ft-backend
ExecStart=/OPSFLEET_ROOT/bin/opsfleet-backend
Restart=on-failure
RestartSec=5
Environment=GIN_MODE=release
# 可选：K8s 制品镜像 manifest（见 deploy/k8s-mirror/README.md）；不存在则忽略
EnvironmentFile=-/etc/opsfleet/backend.env

[Install]
WantedBy=multi-user.target
