# deploy/k8s-mirror/k8s-mirror-ttl-cleanup.timer
[Unit]
Description=Weekly TTL cleanup for OpsFleet K8s mirror

[Timer]
OnCalendar=Sun *-*-* 03:30:00
Persistent=true

[Install]
WantedBy=timers.target
